# configs/raspberrypi/home/tygoe/.config/systemd/user/zrok-dynmap.service
[Unit]
Description=Start sharing mcdash over zrok
After=network.target

[Service]
WorkingDirectory=/home/tygoe/Documents/mcserver
ExecStart=/usr/bin/zrok share reserved dynmap --headless
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
